# systemd/camera-raw.socket
[Socket]
ListenStream=8000
Accept=yes
MaxConnections=1
