# systemd/camera-tee.service
[Service]
ExecStart=/bin/sh -c 'ncat localhost 8001 | ncat -l --send-only -k 8002'
